# DEBIAN/prerm.in
#!/bin/bash


# Other prerm actions
rm_ldconfig() {
  # left-hand term originates from ENABLE_LDCONFIG = ON/OFF at package build
  if [ "@ENABLE_LDCONFIG@" == "ON" ]; then
    rm -f /etc/ld.so.conf.d/x86_64-libamd_smi_lib.conf
    ldconfig
  fi
}


rm_leftovers() {
  # remove pyc files generated by python
  rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/@CMAKE_INSTALL_LIBEXECDIR@/amdsmi_cli/__pycache__"
  rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/@SHARE_INSTALL_PREFIX@/amdsmi/__pycache__"

  # remove build and egg files
  rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/@SHARE_INSTALL_PREFIX@/amdsmi.egg-info"
  rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/@SHARE_INSTALL_PREFIX@/build"

  # remove leftover doc files
  if test -e "@CPACK_PACKAGING_INSTALL_PREFIX@/@SHARE_INSTALL_PREFIX@/../doc/amd_smi*"; then
    rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/@SHARE_INSTALL_PREFIX@/../doc/amd_smi*"
  fi
}

rm_logFolder() {
  rm -rf /var/log/amd_smi_lib
}


rm_rocm_tests_dir(){
  if [ -d "@CPACK_PACKAGING_INSTALL_PREFIX@/share/amd_smi/tests/" ]; then
    rm -rf "@CPACK_PACKAGING_INSTALL_PREFIX@/share/amd_smi/tests/"
    echo "Removed ROCm tests directory."
  fi
}


return_logrotateToOrigConfig() {
  local logrotateConfFile=/etc/logrotate.d/amd_smi.conf
  if [ -f $logrotateConfFile ]; then
    rm -rf "$logrotateConfFile"
  fi
  if [ -f /etc/cron.hourly/logrotate ]; then
    mv /etc/cron.hourly/logrotate /etc/cron.daily/logrotate
  fi
  if [ -f /lib/systemd/system/logrotate.timer.backup ]; then
    cp /lib/systemd/system/logrotate.timer.backup /lib/systemd/system/logrotate.timer
    rm -rf /lib/systemd/system/logrotate.timer.backup
    systemctl reenable --now logrotate.timer
  fi
}

rm_python_lib() {
  # get python version
  local python3_minor_version
  python3_minor_version=$(python3 -c 'import sys;print(sys.version_info.minor)')
  if [ $? -ne 0 ]; then
    echo "[WARNING] Could not determine python version. "\
      "AMD-SMI python library will not be uninstalled."
    return
  fi

  # check if python version is supported
  if [ "$python3_minor_version" -lt 6 ]; then
    echo "[WARNING] AMD-SMI python library is not supported on python version 3.$python3_minor_version. "\
      "AMD-SMI python library will not be uninstalled."
    return
  fi

  # Remove old python library
  local pip_list_output
  pip_list_output=$(python3 -m pip list --format=columns --disable-pip-version-check)
  # check pip list output for amdsmi
  if [[ $pip_list_output == *"amdsmi"* ]]; then
    PIP_ROOT_USER_ACTION=ignore PIP_BREAK_SYSTEM_PACKAGES=1 python3 -m pip uninstall amdsmi --yes --quiet --disable-pip-version-check
  fi

  pip_list_output=$(python3 -m pip list --format=columns --disable-pip-version-check)
  # check pip list output for amdsmi
  if [[ $pip_list_output == *"amdsmi"* ]]; then
    echo "[WARNING] AMD-SMI python library (amdsmi) is still installed in pip. "\
      "Check post install to ensure version is correct"
  else
    echo "Removed AMD-SMI python library (amdsmi)..."
  fi
}


case "$1" in
  ( remove | upgrade)
    # remove old gpuv-smi symlink
    rm -f @CPACK_PACKAGING_INSTALL_PREFIX@/bin/gpuv-smi &> /dev/null
    rm_python_lib
    rm_ldconfig
    rm_leftovers
    rm_logFolder
    rm_rocm_tests_dir
    return_logrotateToOrigConfig
  ;;
  ( purge )
  ;;
  ( * )
    exit 0
  ;;
esac
