# DEBIAN/prerm.in
#!/bin/bash
set -e

packageName="amd-smi-lib"
logPath=/var/log/amd_smi_lib
logName=AMD-SMI-lib.log
logFile="${logPath}/${logName}"
logrotateConfFile=/etc/logrotate.d/amd_smi.conf

rm_logFolder() {
  sudo rm -rf "$logPath"
}

return_logrotateToOrigConfig() {
  if [ -f $logrotateConfFile ]; then
    sudo rm -rf "${logrotateConfFile}"
  fi
  if [ -f /etc/cron.hourly/logrotate ]; then
    sudo mv /etc/cron.hourly/logrotate /etc/cron.daily/logrotate
  fi
  if [ -f /lib/systemd/system/logrotate.timer.backup ]; then
    sudo cp /lib/systemd/system/logrotate.timer.backup /lib/systemd/system/logrotate.timer
    sudo rm -rf /lib/systemd/system/logrotate.timer.backup
    sudo systemctl reenable --now logrotate.timer
  fi
}

rm_ldconfig() {
  # left-hand term originates from ENABLE_LDCONFIG = ON/OFF at package build
  if [ "@ENABLE_LDCONFIG@" == "ON" ]; then
    rm -f /etc/ld.so.conf.d/x86_64-libamd_smi_lib.conf
    ldconfig
  fi
}

rm_pyc() {
 # remove pyc file generated by python
 rm -rf @CPACK_PACKAGING_INSTALL_PREFIX@/@CMAKE_INSTALL_LIBEXECDIR@/amd_smi/__pycache__
}


case "$1" in
  ( remove | upgrade)
    rm_ldconfig
    rm_pyc
    rm_logFolder
    return_logrotateToOrigConfig
  ;;
  ( purge )
  ;;
  ( * )
    exit 0
  ;;
esac
